# roles/shutdown/templates/systemd.service.j2
[Unit]
Description={{ shutdown_timer_description }}

[Service]
Type={{ shutdown_timer_svc_type }}
ExecStart={{ shutdown_timer_svc_cmd }}

[Install]
WantedBy=multi-user.target